# system-cfg-files/files.initrd.d/rhel/configs.list
/etc/rc.d/init.d/functions
/etc/bashrc                          optional
/etc/profile.d/glib2.csh             optional
/etc/profile.d/glib2.sh              optional
/etc/profile.d/colorls.csh           optional
/etc/profile.d/colorls.sh            optional
/etc/profile.d/less.csh              optional
/etc/profile.d/less.sh               optional
/etc/profile.d/which-2.sh            optional
/etc/profile.d/lang.csh              optional
/etc/profile.d/lang.sh               optional
/etc/redhat-release                  optional
/usr/share/hwdata/                   optional
