# system-cfg-files/gfs61-es40-files.x86_64.list
#****h* system-cfg-files/gfs61-es40-files.x86_64.list
#  NAME
#    gfs61-es40-files.x86_64.list
#    $id$
#  DESCRIPTION
#*******
# $Id: gfs61-es40-files.x86_64.list,v 1.6 2006-06-07 09:42:23 marc Exp $
bash
nash
nslookup
head
tail
chmod
df
dd
e2label
grep
expr
sed
echo
hostname
domainname
cat
mount
mv
umount
rm
tee
ifenslave
sort
reboot
rmmod
mkdir
mknod
chroot
fuser
gzip
gunzip
ps
ifconfig
route
insmod
pivot_root
#startproc
strace
/bin/true
lsmod
sleep
killall
lock_gulmd
awk
killall5
sh
sfdisk
init
kudzu
gnbd_import
grubby
modprobe
depmod
ping
dhclient
dhclient-script
ls
uname
basename
consoletype
ifup
ifdown
ip
mii-tool
modinfo
fdisk
cp
date
tar
lspci
pidof
dmesg
find
mv
dirname
ipcalc
arping
mii-tool
fence_ack_manual
fence_gnbd
fence_manual
fence_node
fence_drac
#gfs_conf
gfs_fsck
gfs_grow
gfs_jadd
gfs_mkfs
gfs_quota
gfs_tool
lock_gulmd
mkfs.gfs
mkfs.ext2
grub
grub-install
grub-md5-crypt
mbchk
lilo
xargs
host
which
ldd
/bin/vi
more
ethtool
seq
ln
pwd
wget
gulm_tool
hwclock
tr
syslogd
logger
scsi_info
# START GFS stuff
lvm
#vgscan
#vgchange
udev
udevstart
clvmd
ccs_tool
ccsd
cman_tool
fenced
fence_tool
/etc/cluster/cluster.conf
/etc/modprobe.conf
/sbin/circleping
/sbin/cluster_cmd
/sbin/magma_tool
# END GFS stuff
#
# START iscsi part
# /etc/iscsi.cfg is in template because of emptyness.
iscsid
iscsi-iname
iscsi-mountall
iscsi-device
iscsi-ls
iscsi-umountall
# END iscsi
# START other stuff
# START to be compiled by perlcc
# @perlcc /opt/atix/comoonics_cs/ccs_fileread
# END to be compiled by perlcc
# START special things for lock_gulmd
/opt/atix/comoonics_cs/
/sbin/fence_apc
/sbin/fence_rsa
/sbin/fence_apc_old
/sbin/fence_bladecenter
/sbin/fence_brocade
/sbin/fence_egenera
/sbin/fence_mcdata
/sbin/fence_vixel
/sbin/fence_wti
/sbin/fence_xcat
/sbin/fence_ilo
/usr/bin/perl
# END special things for lock_gulmd
# END other stuff
# START Systemfiles
/dev
/etc/sysconfig/network-scripts/ifdown-aliases
/etc/sysconfig/network-scripts/ifcfg-lo
/etc/sysconfig/network-scripts/ifdown
/etc/sysconfig/network-scripts/ifup-aliases
/etc/sysconfig/network-scripts/ifdown-post
/etc/sysconfig/network-scripts/ifup
/etc/sysconfig/network-scripts/ifup-post
/etc/sysconfig/network-scripts/ifup-routes
/etc/sysconfig/network-scripts/init.ipv6-global
/etc/sysconfig/network-scripts/network-functions
/etc/sysconfig/network-scripts/ifcfg-eth0
/etc/sysconfig/kudzu
/etc/sysconfig/network
/etc/sysconfig/clock
/etc/bashrc
/etc/profile.d/glib2.csh
/etc/profile.d/glib2.sh
/etc/profile.d/colorls.csh
/etc/profile.d/colorls.sh
/etc/profile.d/less.csh
/etc/profile.d/less.sh
/etc/profile.d/which-2.sh
/etc/profile.d/lang.csh
/etc/profile.d/lang.sh
/etc/profile.d/krb5.csh
/etc/profile.d/krb5.sh
/etc/filesystems
/etc/host.conf
/etc/hosts.allow
/etc/hosts.deny
/etc/profile
/etc/protocols
/etc/nsswitch.conf
/etc/ld.so.conf
/etc/localtime
/etc/rc.d
/etc/rc.d/init.d
/etc/rc.d/init.d/functions
/etc/rc.d/init.d/killall
/etc/rc.d/init.d/network
/etc/rc.d/init.d/iscsi
/etc/init.d
/etc/redhat-release
/etc/resolv.conf
/etc/hosts
/etc/iscsi.conf
/etc/initiatorname.iscsi
# Integrated in rpm (can be changed)!!
/usr/share/hwdata/
#/usr/share/hwdata/CardMonitorCombos
#/usr/share/hwdata/Cards
#/usr/share/hwdata/MonitorsDB
#/usr/share/hwdata/pci.ids
#/usr/share/hwdata/upgradelist
#/usr/share/hwdata/usb.ids
#/usr/share/hwdata/pcitable
/usr/lib/locale/en_US/
/lib64/libmagma.so
/lib64/libmagma.so.1
/lib64/libmagma.so.1.0.0
/lib64/libmagma_nt.so
/lib64/libmagma_nt.so.1
/lib64/libmagma_nt.so.1.0.0
/lib64/libmagmamsg.so
/lib64/libmagmamsg.so.1
/lib64/libmagmamsg.so.1.0.0
/lib64/magma
/lib64/magma/magma_gulm.so
/lib64/magma/magma_sm.so
/lib64/libssl.so.0.9.7a
/lib64/libssl.so.0.9.7a
/lib64/libssl.so.4
/lib/ld-2.3.4.so
/lib/ld-linux.so.2
/lib/libBrokenLocale.so.1
/lib/libNoVersion.so.1
/lib/libSegFault.so
/lib/libanl-2.3.4.so
/lib/libanl.so.1
/lib/libc-2.3.4.so
/lib/libc.so.6
/lib/libcrypt-2.3.4.so
/lib/libcrypt.so.1
/lib/libdl-2.3.4.so
/lib/libdl.so.2
/lib/libm-2.3.4.so
/lib/libm.so.6
/lib/libnsl-2.3.4.so
/lib/libnsl.so.1
/lib/libnss1_compat-2.3.4.so
/lib/libnss1_compat.so.1
/lib/libnss1_dns-2.3.4.so
/lib/libnss1_dns.so.1
/lib/libnss1_files-2.3.4.so
/lib/libnss1_files.so.1
/lib/libnss1_nis-2.3.4.so
/lib/libnss1_nis.so.1
/lib/libnss_compat-2.3.4.so
/lib/libnss_compat.so.1
/lib/libnss_compat.so.2
/lib/libnss_dns-2.3.4.so
/lib/libnss_dns.so.1
/lib/libnss_dns.so.2
/lib/libnss_files-2.3.4.so
/lib/libnss_files.so.1
/lib/libnss_files.so.2
/lib/libnss_hesiod-2.3.4.so
/lib/libnss_hesiod.so.2
/lib/libnss_nis-2.3.4.so
/lib/libnss_nis.so.1
/lib/libnss_nis.so.2
/lib/libnss_nisplus-2.3.4.so
/lib/libnss_nisplus.so.2
/lib/libpthread-0.10.so
/lib/libpthread.so.0
/lib/libresolv-2.3.4.so
/lib/libresolv.so.2
/lib/librt-2.3.4.so
/lib/librt.so.1
/lib/libthread_db-1.0.so
/lib/libthread_db.so.1
/lib/libutil-2.3.4.so
/lib/libutil.so.1
/lib64/libtermcap.so.2
/lib64/libtermcap.so.2.0.8
/lib64/libssl.so.4
/lib64/libssl.so.0.9.7a
/lib64/libcrypto.so.0.9.7a
/lib64/libcrypto.so.4
# $Log: gfs61-es40-files.x86_64.list,v $
# Revision 1.6  2006-06-07 09:42:23  marc
# *** empty log message ***
#
# Revision 1.5  2006/05/12 13:03:55  marc
# Some new files for Version 1.0
#
# Revision 1.4  2006/05/07 11:37:45  marc
# added tail and head
#
# Revision 1.3  2006/05/03 12:47:48  marc
# added documentation
#
# Revision 1.2  2006/04/13 18:50:59  marc
# readded /usr/share/hwdata
#
# Revision 1.1  2006/04/11 13:42:26  marc
# initial revision
#
# Revision 1.1  2006/04/09 16:31:27  marc
# initial revision
#
# Revision 1.2  2005/07/08 12:59:56  mark
# added some files
#
# Revision 1.2  2005/01/03 09:08:05  marc
# added tags
#
/usr/lib64/libwrap.so.0.7.6
/usr/lib64/libwrap.so.0
/usr/lib/perl5/
/usr/lib64/perl5/
/usr/lib64/libexpat.so.0
/usr/lib64/libz.so.1
/usr/lib64/libz.so.1.2.1.2
/usr/lib64/libdes425.so.3
/usr/lib64/libdes425.so.3.0
/usr/lib64/libgssapi_krb5.so.2
/usr/lib64/libgssapi_krb5.so.2.2
/usr/lib64/libgssrpc.so.3
/usr/lib64/libgssrpc.so.3.0
/usr/lib64/libk5crypto.so.3
/usr/lib64/libk5crypto.so.3.0
/usr/lib64/libkadm5clnt.so.5
/usr/lib64/libkadm5clnt.so.5.1
/usr/lib64/libkadm5srv.so.5
/usr/lib64/libkadm5srv.so.5.1
/usr/lib64/libkdb5.so.4
/usr/lib64/libkdb5.so.4.0
/usr/lib64/libkrb4.so.2
/usr/lib64/libkrb4.so.2.0
/usr/lib64/libkrb5.so.3
/usr/lib/locale/en_US/
/usr/lib/locale/en_US/
/usr/lib/locale/en_US/
/lib64/libnss_compat-2.3.4.so
/lib64/libnss_compat.so.2
/lib64/libnss_db.so.2
/lib64/libnss_db.so.2.0.0
/lib64/libnss_dns-2.3.4.so
/lib64/libnss_dns.so.2
/lib64/libnss_files-2.3.4.so
/lib64/libnss_files.so.2
/lib64/libnss_hesiod-2.3.4.so
/lib64/libnss_hesiod.so.2
/lib64/libnss_ldap-2.3.4.so
/lib64/libnss_ldap.so.2
/lib64/libnss_nis-2.3.4.so
/lib64/libnss_nisplus-2.3.4.so
/lib64/libnss_nisplus.so.2
/lib64/libnss_nis.so.2
