# system-cfg-files/files.initrd.d/configs.list
/etc/sysconfig/kudzu
/etc/sysconfig/network
/etc/sysconfig/clock
/etc/bashrc
/etc/profile.d/glib2.csh
/etc/profile.d/glib2.sh
/etc/profile.d/colorls.csh
/etc/profile.d/colorls.sh
/etc/profile.d/less.csh
/etc/profile.d/less.sh
/etc/profile.d/which-2.sh
/etc/profile.d/lang.csh
/etc/profile.d/lang.sh
/etc/profile.d/krb5.csh
/etc/profile.d/krb5.sh
/etc/filesystems
/etc/host.conf
/etc/hosts.allow
/etc/hosts.deny
/etc/profile
/etc/protocols
/etc/nsswitch.conf
/etc/ld.so.conf
/etc/localtime
/etc/rc.d
/etc/rc.d/init.d
/etc/rc.d/init.d/functions
/etc/rc.d/init.d/killall
/etc/rc.d/init.d/network
/etc/init.d
/etc/redhat-release
/etc/resolv.conf
/etc/hosts
/etc/modprobe.conf
/usr/share/hwdata/
