# system-cfg-files/files.initrd.d/configs.list
/etc/sysconfig/network
/etc/filesystems
/etc/host.conf
/etc/hosts.allow
/etc/hosts.deny
/etc/profile
/etc/protocols
/etc/nsswitch.conf
/etc/ld.so.conf
/etc/localtime
/etc/rc.d
/etc/rc.d/init.d
/etc/init.d
/etc/resolv.conf
/etc/hosts
/etc/modprobe.conf
